# scripts/05_automation/systemd/buckup.timer
[Unit]
Description=Timer para ejecutar backup.sh
[Timer]
OnCalendar=Fry 15:00
Persistent=true
[Install]
WantedBy=timers.target
